# fedora-kafka/assets/supervisord.conf
[supervisord]
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:kafka]
command=/bin/kafka-start
exitcodes=0,2
environment=KAFKA_HOME=%(ENV_KAFKA_HOME)s,KAFKA_RUN=%(ENV_KAFKA_RUN)s,KAFKA_OPTS="%(ENV_KAFKA_OPTS)s"
