# fedora-zookeeper/assets/supervisord.conf
[supervisord]
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:zookeeper]
command=/bin/zookeeper-start
environment=ZOO_CONF=%(ENV_ZOO_CONF)s
