# fedora-postgres/assets/supervisord.conf
[supervisord]
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:postgresql]
command=/bin/pgsql-server-start
user=postgres
exitcodes=0,2
environment=PGHOME=%(ENV_PGHOME)s,PGDATA=%(ENV_PGDATA)s,PGPORT=%(ENV_PGPORT)s
