# fedora-activemq/assets/supervisord.conf
[supervisord]
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:activemq]
command=/bin/activemq-start
exitcodes=0,2
environment=ACTIVEMQ_HOME=%(ENV_ACTIVEMQ_HOME)s,ACTIVEMQ_RUN=%(ENV_ACTIVEMQ_RUN)s
